# lib/capistrano-rails-server/recipes/templates/nginx_init.erb
#! /bin/sh

### BEGIN INIT INFO
# Provides:          nginx
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts the nginx web server
# Description:       starts nginx using start-stop-daemon
### END INIT INFO

PATH=/opt/nginx/sbin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=/opt/nginx/sbin/nginx
NAME=nginx
DESC=nginx
PIDFILE=/opt/nginx/logs/$NAME.pid

test -x $DAEMON || exit 0

# Include nginx defaults if available
if [ -f /etc/default/nginx ] ; then
        . /etc/default/nginx
fi

set -e

case "$1" in
  start)
        echo -n "Starting $DESC: "
        if [ ! -f $PIDFILE ]
        then
            start-stop-daemon --start --quiet --pidfile $PIDFILE --exec $DAEMON -- $DAEMON_OPTS
            echo "$NAME."
        else
            echo "Already running."
        fi
        ;;
  stop)
        echo -n "Stopping $DESC: "
        if [ -f $PIDFILE ]
        then
            start-stop-daemon --stop --quiet --pidfile $PIDFILE --exec $DAEMON
            echo "$NAME."
        else
            echo "Already not running."
        fi
        ;;
  restart|force-reload)
        echo -n "Restarting $DESC: "
        if [ -f $PIDFILE ]
        then
            start-stop-daemon --stop --quiet --pidfile $PIDFILE --exec $DAEMON
            sleep 1
        fi
        start-stop-daemon --start --quiet --pidfile $PIDFILE --exec $DAEMON -- $DAEMON_OPTS
        echo "$NAME."
        ;;
  reload)
          echo -n "Reloading $DESC configuration: "
          if [ -f $PIDFILE ]
          then
              start-stop-daemon --stop --signal HUP --quiet --pidfile $PIDFILE --exec $DAEMON
              echo "$NAME."
          else
              echo "Not running, can't reload."
          fi
          ;;
      *)
            N=/etc/init.d/$NAME
            echo "Usage: $N {start|stop|restart|reload|force-reload}" >&2
            exit 1
            ;;
    esac
    exit 0
